# templates/default/systemd/docker.socket.erb
[Unit]
Description=Docker Socket for the API
PartOf=docker.service

[Socket]
<% Array(@config.host).each do |socket| -%>
ListenStream=<%= socket.gsub(/^unix:\/\//,'').gsub(/^tcp:\/\//,'').gsub(/^fd:\/\//,'') %>
<% end -%>
SocketMode=0660
SocketUser=root
SocketGroup=<%= @config.group || 'root' %>

[Install]
WantedBy=sockets.target
